# dashboard.conf
[program:dashboard]
command = /home/dashboard/venv/bin/gunicorn --workers 3 -b 0.0.0.0:8000 dashboard.wsgi
directory = /home/dashboard/dashboard
user = dashboard
autostart = true
autorestart = true
stderr_logfile = /var/log/dashboard.err.log
stdout_logfile = /var/log/dashboard.out.log

